# docker/resources/docker/services.d/whisperfusion/run
#!/bin/bash
cd /root/
./run-whisperfusion.sh
